# assets/supervisor/conf.d/mysql.conf
[program:mysql]
priority = 1
command = service mysql start
autostart = true
user = root
redirect_stderr=true
stdout_logfile = /var/log/supervisor/mysql.log
loglevel = warn
